# tilt.supervisor.conf
[program:node-tilt]
directory=/opt/node-tilt
command=npm start --production
autostart=true
autorestart=true
stderr_logfile=/opt/node-tilt/error.log
user=morken